# docker/production/php-fpm/supervisord.ini
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
user=root
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:schedule-run]
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/storage/logs/schedule.log
stderr_logfile=/dev/stderr

[program:queue-worker]
command=php /var/www/artisan queue:work --timeout=600
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/storage/logs/queue.log
stderr_logfile=/dev/stderr
